# configs/stockanalysis-resource-monitor.service
[Unit]
Description=StockAnalysis Resource Monitor
After=redis-server.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hacker
WorkingDirectory=/home/hacker/StockAnalysis
Environment=PYTHONPATH=/home/hacker/StockAnalysis
Environment=REDIS_URL=redis://localhost:6379
ExecStart=/home/hacker/StockAnalysis/.venv/bin/python services/resource_monitor/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
CPUQuota=10%
MemoryMax=100M

[Install]
WantedBy=multi-user.target
