# configs/stockanalysis-notification.service
[Unit]
Description=StockAnalysis Notification Service
After=redis-server.service network-online.target
Wants=network-online.target
Requires=redis-server.service

[Service]
Type=simple
User=hacker
WorkingDirectory=/home/hacker/StockAnalysis
Environment=PYTHONPATH=/home/hacker/StockAnalysis
Environment=REDIS_URL=redis://localhost:6379
ExecStart=/home/hacker/StockAnalysis/.venv/bin/python services/notification-service/main.py --consumer-name prod-1
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
CPUQuota=40%
MemoryMax=300M

[Install]
WantedBy=multi-user.target
